# ops/roles/cpetersen.registrator/templates/registrator.conf
description "Registrator"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

script
  BIND=`ifconfig eth0 | grep "inet addr" | awk '{ print substr($2,6) }'`
  exec DOCKER_HOST=http://$BIND:4243 /opt/gocode/bin/registrator consul:
end script
