# ops/roles/cpetersen.registrator/templates/registrator.conf
description "Registrator"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

script
  exec DOCKER_HOST=localhost:4243 /opt/gocode/bin/registrator consul:
end script
